# config/systemd/watchlist.timer
[Unit]
Description=Run Ross DIY Watchlist each weekday

[Timer]
OnCalendar=Mon..Fri 14:15
Persistent=true

[Install]
WantedBy=timers.target
